# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:dbus-daemon]
command=/usr/bin/dbus-daemon --nofork --nopidfile --session --address=unix:path=/run/dbus/session_bus_socket
autostart=true
autorestart=true
priority=10
environment=DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/session_bus_socket"

[program:signal-cli]
command=/usr/bin/signal-cli daemon --dbus
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/signal-cli.log
stderr_logfile=/var/log/supervisor/signal-cli-error.log
environment=HOME="/var/lib/signal-cli",DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/session_bus_socket"

[program:swb-bridge]
command=/app/.venv/bin/python -m swb
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/swb-bridge.log
stderr_logfile=/var/log/supervisor/swb-bridge-error.log
environment=SIGNAL_WS_HOST="%(ENV_SIGNAL_WS_HOST)s",SIGNAL_WS_PORT="%(ENV_SIGNAL_WS_PORT)s",SIGNAL_DBUS_BUS="%(ENV_SIGNAL_DBUS_BUS)s",SIGNAL_LOG_LEVEL="%(ENV_SIGNAL_LOG_LEVEL)s",SIGNAL_WS_TOKEN="%(ENV_SIGNAL_WS_TOKEN)s",SIGNAL_ACCOUNT="%(ENV_SIGNAL_ACCOUNT)s",DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/session_bus_socket"
